# configs/debian/sysmon-config.xml
<Sysmon schemaversion="4.81">
  <HashAlgorithms>sha256</HashAlgorithms>
  <CheckRevocation/>

  <EventFiltering> 
    <RuleGroup name="Process_Creation_Events" groupRelation="or">
    <!-- we nare going to exclude logging for paths that are noisy for event id 1 -->
            <ProcessCreate onmatch="exclude">
                <Image condition="is">/var/ossec/bin/wazuh-agent</Image>
                <Image condition="is">/lib/systemd/systemd</Image>
                <Image condition="is">/usr/sbin/rsyslogd</Image>
                <Image condition="is">/usr/bin/dbus-daemon</Image>
                <Image condition="is">/lib/systemd/systemd-journald</Image>
                <Image condition="is">/usr/bin/apt</Image>
                <Image condition="is">/var/bin/dpkg</Image>
                <Image condition="is">/lib/systemd/systemd-udevd</Image>
                <Image condition="is">/usr/bin/mandb</Image>
            </ProcessCreate>

            <ProcessCreate onmatch="include">
                <!-- shells -->
                <Image condition="is">/bin/bash</Image>
                <Image condition="is">/bin/sh</Image>
                <Image condition="is">/bin/dash</Image>
                <Image condition="is">/usr/local/bin/zsh</Image>
                <!-- network tools -->
                <Image condition="end with">/nmap</Image>
                <Image condition="end with">/curl</Image>
                <Image condition="end with">/wget</Image>
                <Image condition="end with">/nc</Image>
                <Image condition="end with">/netcat</Image>
                <Image condition="end with">/socat</Image>

                <!-- discovery -->
                <Image condition="end with">/whoami</Image>
                <Image condition="end with">/id</Image>
                <Image condition="end with">/uname</Image>
                <Image condition="end with">/ifconfig</Image>
                <Image condition="end with">/ip</Image>

                <!-- priv escalation -->
                <Image condition="end with">/sudo</Image>
                <Image condition="end with">/su</Image>
                <Image condition="end with">/pkexec</Image>

                <!-- bypass firewalls or decode paylods -->
                <Image condition="end with">/base64</Image>
                <Image condition="end with">/openssl</Image>
                <Image condition="contains">/python</Image> <!--check different versions of python -->
                <Image condition="end with">/perl</Image>

                <!--Persistence -->
                <Image condition="end with">/useradd</Image>
                <Image condition="end with">/systemctl</Image>
                <Image condition="end with">/crontab</Image>
                <Image condition="end with">/adduser</Image>
            </ProcessCreate>
    </RuleGroup>


    <RuleGroup name="T1059.004-Unix-Shell" groupRelation="or">
        <ProcessCreate onmatch="include">
                <Rule groupRelation="and">
                        <ParentImage condition="contains">/curl</ParentImage>
                        <Image condition="end with">/sh</Image> 
                </Rule>

                <Rule groupRelation="and">
                        <User condition="is">root</User>
                        <Image condition="end with">/bash</Image>
                </Rule>
        </ProcessCreate>
    </RuleGroup>

<!--https://attack.mitre.org/detectionstrategies/DET0290/ technic T1053.003 -->    
    
    <RuleGroup name="Cron_Persistence_Monitoring" groupRelation="or">
      <FileCreate onmatch="include">
        <TargetFilename condition="begin with">/etc/cron.</TargetFilename>
        <TargetFilename condition="begin with">/var/spool/cron/crontabs/</TargetFilename>
        <TargetFilename condition="is">/etc/crontab</TargetFilename>
      </FileCreate>
      
      <ProcessCreate onmatch="include">
        <Image condition="is">/usr/bin/crontab</Image>
      </ProcessCreate>
    </RuleGroup>
    
    <RuleGroup name="Cron_Execution_Anomalies" groupRelation="or">
        <ProcessCreate onmatch="include">
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">base64</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">python</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">perl</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">bash -i</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">curl</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">wget</CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">nc </CommandLine>
                </Rule>
                <Rule groupRelation="and">
                        <ParentImage condition="contains">cron</ParentImage>
                        <CommandLine condition="contains">/dev/tcp/</CommandLine>
                </Rule>
        </ProcessCreate>
   </RuleGroup>


<!-- end of the rules for T1053.003 --> 

<!--detetcion rule DET0110 for sub technic T1548.001 -->
    <RuleGroup name="Privilege_Escalation_Setuid_Setgid" groupRelation="or">
      
      <ProcessCreate onmatch="include">
        <CommandLine condition="contains">4755</CommandLine>
        <CommandLine condition="contains">2755</CommandLine>
        <CommandLine condition="contains">4000</CommandLine>
        <CommandLine condition="contains">u+s</CommandLine>
        <CommandLine condition="contains">g+s</CommandLine>
      </ProcessCreate>

      <ProcessCreate onmatch="include">
        <Rule groupRelation="and">
          <Image condition="begin with">/tmp/</Image>
          <User condition="is">root</User>
        </Rule>
        <Rule groupRelation="and">
          <Image condition="begin with">/var/tmp/</Image>
          <User condition="is">root</User>
        </Rule>
        <Rule groupRelation="and">
          <Image condition="begin with">/dev/shm/</Image>
          <User condition="is">root</User>
        </Rule>
      </ProcessCreate>

    </RuleGroup>

<!--end of T1548.001 -->

<!-- T1548.003  -->
    <RuleGroup name="Sudo_Privilege_Abuse" groupRelation="or">
        <ProcessCreate onmatch="include">
                <Rule groupRelation="and">
                        <ParentImage condition="end with">/sudo</ParentImage>
                        <Image condition="contains any">/bin/bash;/bin/sh;/bin/dash;/usr/bin/python3</Image>
                </Rule>
    
                <Rule groupRelation="and">
                        <ParentImage condition="end with">/sudo</ParentImage>
                        <Image condition="end with">/visudo</Image>
                </Rule>

                <Rule groupRelation="and">
                        <CommandLine condition="contains">/etc/sudoers</CommandLine>
                        <Image condition="contains any">nano;vi;vim;sed;echo;tee</Image>
                </Rule>
        </ProcessCreate>

        <FileCreate onmatch="include">
                <TargetFilename condition="begin with">/etc/sudoers</TargetFilename>
        </FileCreate>
    </RuleGroup>    

<!-- end of T1548.003  -->

    <RuleGroup name="Network_Connection_Events" groupRelation="or">
        <NetworkConnect onmatch="exclude">
            <Image condition="is">/usr/sbin/rsyslogd</Image>
            <Image condition="is">/var/ossec/bin/wazuh-agent</Image>
            <Image condition="is">/usr/bin/dbus-daemon</Image>
            <Image condition="is">/lib/systemd/systemd-resolved</Image>
            <Image condition="is">/usr/lib/apt/methods/http</Image>
            <Image condition="is">/usr/lib/apt/methods/https</Image>
            <DestinationIp condition="is">127.0.0.1</DestinationIp>
        </NetworkConnect>

        <NetworkConnect onmatch="include">
            <Image condition="end with">/bash</Image>
            <Image condition="end with">/sh</Image>
            <Image condition="end with">/dash</Image>       
            <Image condition="end with">/curl</Image>
            <Image condition="end with">/wget</Image>
            <Image condition="end with">/nc</Image>
            <Image condition="end with">/netcat</Image>
            <Image condition="end with">/socat</Image>
            <Image condition="end with">/nmap</Image>
            <Image condition="end with">/scp</Image>
            <Image condition="end with">/sftp</Image>
            <Image condition="contains">/python</Image>
            <Image condition="end with">/perl</Image>
            <Image condition="end with">/ruby</Image>
            <Image condition="end with">/lua</Image>
            <Image condition="end with">/node</Image>
            <Image condition="is">/usr/bin/ssh</Image>
        </NetworkConnect>
    </RuleGroup>     


    <RuleGroup name="Process_Termination_Events" groupRelation="or">
      <ProcessTerminate onmatch="exclude">
        <Image condition="is">/usr/bin/clear</Image>
        <Image condition="is">/usr/bin/grep</Image>
        <Image condition="is">/usr/bin/sed</Image>
        <Image condition="is">/usr/bin/awk</Image>
        <Image condition="is">/usr/bin/cat</Image>
      </ProcessTerminate>

      <ProcessTerminate onmatch="include">
        <Image condition="end with">/bash</Image>
        <Image condition="end with">/sh</Image>
        <Image condition="end with">/dash</Image>

        <Image condition="is">/var/ossec/bin/wazuh-agent</Image>
        <Image condition="is">/usr/sbin/rsyslogd</Image>
        <Image condition="is">/usr/bin/sysmon</Image>
        
        <Image condition="end with">/nc</Image>
        <Image condition="end with">/netcat</Image>
        <Image condition="end with">/socat</Image>
        <Image condition="contains">/python</Image>
      </ProcessTerminate>
    </RuleGroup>


    <RuleGroup name="File_Creation_Events" groupRelation="or">
      <FileCreate onmatch="exclude">
        <TargetFilename condition="begin with">/var/log</TargetFilename>
        <TargetFilename condition="begin with">/tmp/sh-thd</TargetFilename>
        <TargetFilename condition="contains">/var/lib/apt/lists</TargetFilename>
        <TargetFilename condition="begin with">/var/cache</TargetFilename>
      </FileCreate>

      <FileCreate onmatch="include">
        <TargetFilename condition="contains">.ssh/authorized_keys</TargetFilename>
        <TargetFilename condition="contains">.ssh/id_rsa</TargetFilename>

        <TargetFilename condition="begin with">/etc/cron</TargetFilename>
        <TargetFilename condition="begin with">/var/spool/cron/crontabs</TargetFilename>
        <TargetFilename condition="begin with">/etc/systemd/system</TargetFilename>

        <TargetFilename condition="begin with">/var/www/html</TargetFilename>

        <TargetFilename condition="begin with">/usr/local/bin</TargetFilename>
        <TargetFilename condition="begin with">/usr/bin</TargetFilename>
        <TargetFilename condition="begin with">/bin</TargetFilename>

        <TargetFilename condition="is">/etc/passwd</TargetFilename>
        <TargetFilename condition="is">/etc/sudoers</TargetFilename>
        <TargetFilename condition="begin with">/etc/sudoers.d</TargetFilename>
        <TargetFilename condition="is">/etc/shadow</TargetFilename>
        <TargetFilename condition="is">/etc/ld.so.preload</TargetFilename> </FileCreate>
    </RuleGroup>    


  </EventFiltering>
</Sysmon>
